# roles/user/files/dotfiles/.config/autostart/ignore-lid-switch-tweak.desktop
[Desktop Entry]
Type=Application
Name=ignore-lid-switch-tweak
Exec=/bin/sh -c '$HOME/.local/lib/gnome-tweak-tool-lid-inhibitor'
