# roles/user/files/dotfiles/.config/systemd/user/gnome-terminal-server.service.d/override.conf
[Service]
# workaround for:
# https://gitlab.gnome.org/GNOME/vte/-/issues/347
# https://github.com/googlefonts/noto-cjk/issues/144#issuecomment-478968910
Environment=LC_CTYPE=C.UTF-8
