# roles/laptop_t14gen3amd/files/t14-gen3-trackpoint-reset.service
[Unit]
Description=Reset Trackpoint

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'rmmod psmouse; modprobe psmouse'
