# roles/laptop_t14gen3amd/files/t14-gen3-trackpoint-reset.timer
[Unit]
Description=Reset Trackpoint periodically

[Timer]
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
